# debian-vscode/vscode-wayland.service
[Unit]
Description=VS Code in Cage (Wayland)
After=systemd-user-sessions.service
Conflicts=getty@tty1.service

[Service]
User=kiosk
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=tty
StandardError=journal
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/local/bin/vscode-wayland.sh
Restart=on-failure
RestartSec=2
KillMode=process
WorkingDirectory=/home/kiosk
SupplementaryGroups=tty video input render
ExecStartPre=+/bin/sh -c 'mkdir -p /run/user/1000 && chown kiosk:kiosk /run/user/1000 && chmod 700 /run/user/1000'

[Install]
WantedBy=multi-user.target